# systemd/tsocial.service
[Unit]
Description=terminal-social
After=network.target

[Service]
User=deploy
Group=deploy
ExecStart=/usr/local/bin/tsocial
Restart=on-failure
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tsocial
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
Environment=TSOCIAL_LISTEN=0.0.0.0:2222
Environment=TSOCIAL_DB=/var/lib/tsocial/tsocial.db
Environment=TSOCIAL_HOST_KEY=/var/lib/tsocial/host_key

[Install]
WantedBy=multi-user.target
